# dind-cache/dind-wipe.service
[Unit]
Description=DIND Cleanup Job

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart dind-cache.service
